# itsm-chat/deploy/itsm-chat.service
[Unit]
Description=ITSM Chat Backend API
After=network.target

[Service]
Type=simple
User=paperclip
WorkingDirectory=/home/paperclip/app/itsm-chat/backend
Environment=PATH=/home/paperclip/app/itsm-chat/backend/venv/bin:/usr/bin
ExecStart=/home/paperclip/app/itsm-chat/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8100
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
